# scripts/com.portfolio-viz.sync-sheets.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.portfolio-viz.sync-sheets</string>

	<key>ProgramArguments</key>
	<array>
		<string>PROJECT_ROOT/.venv/bin/python</string>
		<string>PROJECT_ROOT/scripts/sync_to_sheets.py</string>
	</array>

	<key>WorkingDirectory</key>
	<string>PROJECT_ROOT</string>

	<!-- 4:15 PM ET, Monday–Friday (Weekday 1=Mon … 5=Fri) -->
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key><integer>1</integer>
			<key>Hour</key><integer>16</integer>
			<key>Minute</key><integer>15</integer>
		</dict>
		<dict>
			<key>Weekday</key><integer>2</integer>
			<key>Hour</key><integer>16</integer>
			<key>Minute</key><integer>15</integer>
		</dict>
		<dict>
			<key>Weekday</key><integer>3</integer>
			<key>Hour</key><integer>16</integer>
			<key>Minute</key><integer>15</integer>
		</dict>
		<dict>
			<key>Weekday</key><integer>4</integer>
			<key>Hour</key><integer>16</integer>
			<key>Minute</key><integer>15</integer>
		</dict>
		<dict>
			<key>Weekday</key><integer>5</integer>
			<key>Hour</key><integer>16</integer>
			<key>Minute</key><integer>15</integer>
		</dict>
	</array>

	<key>StandardOutPath</key>
	<string>PROJECT_ROOT/logs/sync-sheets.log</string>

	<key>StandardErrorPath</key>
	<string>PROJECT_ROOT/logs/sync-sheets.error.log</string>

	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
